# deploy/docker/rootfs/etc/s6-overlay/s6-rc.d/backend/finish
#!/bin/sh
# Crash detection: bring the whole container down if the backend dies rather than
# let s6 restart it in place. See halt-container.sh for the full rationale.
exec /etc/s6-overlay/scripts/halt-container.sh backend "$@"
