# deploy/docker/rootfs/etc/s6-overlay/s6-rc.d/nginx/finish
#!/bin/sh
# Crash detection: bring the whole container down if nginx dies rather than let
# s6 restart it in place. See halt-container.sh for the full rationale.
exec /etc/s6-overlay/scripts/halt-container.sh nginx "$@"
